# etc/cron.d/fun-stuff
# Lets have fortune and cowsay create an motd file - runs every hour and at reboot
@reboot root /usr/local/bin/update-motd &>/dev/null
05 * * * * root /usr/local/bin/update-motd &>/dev/null
