# etc/cron.d/process-check
# Check processes and restart if
*/15 * * * * root [ -f "/root/.local/bin/process-check.sh" ] && /root/.local/bin/process-check.sh &>/dev/null
