# etc/cron.d/process-check
# Check processes and restart if
*/15 * * * * root process-check &>/dev/null
