# etc/cron.d/yum-update
# run yum update - runs daily at 2 am
0 2 * * * root ping -c 1 google.com >/dev/null 2>&1 && yum update -y >/dev/null 2>&1
